# disk-watchdog.service
[Unit]
Description=Adaptive Disk Space Watchdog
Documentation=https://github.com/radrob2/disk-watchdog
After=network.target local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/disk-watchdog
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/lib/disk-watchdog /run
PrivateTmp=true
NoNewPrivileges=true

# Note: No WatchdogSec because main loop has adaptive sleep intervals (1s-300s)

[Install]
WantedBy=multi-user.target
